# rpm/etop.service
[Unit]
Description=etop system monitor for linux system

[Service]
Type=simple
ExecStart=/usr/bin/etop record

[Install]
WantedBy=multi-user.target